# discord/systemd/pr-reviews.service
[Unit]
Description=PR review sweep (alexrudy)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/dev/discord/discord
EnvironmentFile=-%h/dev/discord/discord/.env
# uv ships inside the discord monorepo (not on PATH) — invoke it by
# absolute path. PATH is still set for any subprocess uv-run launches.
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=%h/dev/discord/discord/.local/bin/uv run misc/users/alexrudy/exp-pr-reviews/pr_reviews.py sweep --apply
